# autoload.ini
[program:pepper_main]
command=/usr/bin/python /home/nao/main.py --ip 127.0.0.1 --port 9559
directory=/home/nao
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile=/home/nao/.local/share/PackageManager/logs/pepper_main.log
